# roles/netbootxyz/templates/iso-mount.mount.j2
# {{ ansible_managed }}
# Systemd mount unit for {{ item.item.name }}
# Automatically generated by rhdc.network.netbootxyz role

[Unit]
Description=Loop mount {{ item.item.name }} ISO for HTTP repository access
Documentation=https://github.com/flyemsafe/rhdc-ansible
After=local-fs.target
Requires=local-fs.target
ConditionPathExists={{ item.item.iso_path }}

[Mount]
What={{ item.item.iso_path }}
Where={{ item.item.mount_point }}
Type=iso9660
Options=ro,loop,nofail,context=system_u:object_r:svirt_sandbox_file_t:s0

[Install]
WantedBy=multi-user.target
