# roles/traefik/templates/traefik.container.j2
# Traefik Reverse Proxy - Quadlet Container Unit
# Generated by Ansible - DO NOT EDIT MANUALLY
# Role: rhdc.network.traefik
# Host: {{ traefik_host }}
# Date: {{ ansible_date_time.iso8601 }}

[Unit]
Description=Traefik Reverse Proxy
Documentation=https://doc.traefik.io/traefik/
{% for service in traefik_systemd_after %}
After={{ service }}
{% endfor %}
Wants=network-online.target

[Container]
Image={{ traefik_image }}
ContainerName={{ traefik_container_name }}

# Environment variables
{% if traefik_acme_enabled and traefik_acme_dns_provider == 'cloudflare' %}
Environment=CF_DNS_API_TOKEN={{ traefik_cloudflare_api_token }}
Environment=CF_API_EMAIL={{ traefik_cloudflare_email }}
{% endif %}

# Volume mounts
Volume={{ traefik_config_dir }}:/config:Z
Volume={{ traefik_logs_dir }}:/logs:Z

# Port mappings
PublishPort={{ traefik_http_port }}:{{ traefik_http_port }}
{% if traefik_https_port is defined %}
PublishPort={{ traefik_https_port }}:{{ traefik_https_port }}
{% endif %}
{% if traefik_dashboard_enabled and traefik_dashboard_insecure %}
PublishPort={{ traefik_dashboard_port }}:{{ traefik_dashboard_port }}
{% endif %}
{% if traefik_metrics_prometheus_enabled %}
PublishPort={{ traefik_metrics_prometheus_port }}:{{ traefik_metrics_prometheus_port }}
{% endif %}

# Command arguments
Exec=--configFile=/config/traefik.yml

# Security options
SecurityLabelDisable=false
User={{ traefik_user }}

# Resource limits
{% if traefik_cpu_limit %}
CPUQuota={{ traefik_cpu_limit }}
{% endif %}
{% if traefik_memory_limit %}
Memory={{ traefik_memory_limit }}
{% endif %}

# Health check
{% if traefik_healthcheck_enabled %}
HealthCmd=wget --no-verbose --tries=1 --spider http://localhost:{{ traefik_dashboard_port }}{{ traefik_healthcheck_path }} || exit 1
HealthInterval=30s
HealthRetries=3
HealthStartPeriod=40s
HealthTimeout=10s
{% endif %}

# Restart policy
{% if traefik_auto_restart %}
Restart={{ traefik_systemd_restart_policy }}
RestartSec={{ traefik_restart_sec }}
{% endif %}

# Network
Network=slirp4netns:port_handler=slirp4netns
# Note: slirp4netns enables rootless networking without host network access

[Service]
# Systemd service options
TimeoutStartSec=300
Type=notify
NotifyAccess=all

[Install]
WantedBy=default.target
